# root/etc/services.d/rpcbind/run
#!/usr/bin/with-contenv bash

rpcbind -f
