# root/etc/s6-overlay/s6-rc.d/svc-rpcbind/run
#!/usr/bin/with-contenv bash
set -e

RPCBIND_ARGS="${RPCBIND_ARGS:- -fw}";
echo "Starting rpcbind...";
exec \
    rpcbind ${RPCBIND_ARGS};

