# root/etc/services.d/mountd/run
#!/usr/bin/execlineb -P

foreground { rpc.nfsd -N 2 -N 3 --debug 8  }
rpc.mountd -F -N 2 -N 3 --debug all
